# youtube-web-system.service
[Unit]
Description=Charles's YouTube Web System
After=network.target
Requires=network.target

[Service]
Type=simple
User=jbyrd
Group=jbyrd
WorkingDirectory=/home/jbyrd/hatter-pai-youtube
Environment=PATH=/home/jbyrd/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=youtube-web-system

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/jbyrd/hatter-pai-youtube /mnt/nfs_share/charles/youtube /tmp

[Install]
WantedBy=multi-user.target
