# ansible/templates/grafana-server.service.j2
# Grafana systemd service
# Generated by Ansible - DO NOT EDIT MANUALLY
# Template: grafana-server.service.j2

[Unit]
Description=Grafana - PAI Infrastructure Visualization
Documentation=http://docs.grafana.org
After=network-online.target postgresql.service mariadb.service mysql.service
Wants=network-online.target

[Service]
Type=notify
User={{ grafana_user }}
Group={{ grafana_user }}
RuntimeDirectory=grafana
RuntimeDirectoryMode=0755
ExecStart=/usr/sbin/grafana-server \
  --config={{ grafana_config_dir }}/grafana.ini \
  --pidfile=/var/run/grafana/grafana-server.pid \
  --packaging=rpm \
  cfg:default.paths.logs=/var/log/grafana \
  cfg:default.paths.data={{ grafana_data_dir }} \
  cfg:default.paths.plugins={{ grafana_data_dir }}/plugins \
  cfg:default.paths.provisioning={{ grafana_config_dir }}/provisioning

ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=20
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths={{ grafana_data_dir }}
ReadWritePaths=/var/log/grafana
ReadWritePaths=/var/run/grafana
ReadOnlyPaths={{ grafana_config_dir }}

# Resource limits
MemoryLimit=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
