# ansible/templates/node_exporter.service.j2
# Node Exporter systemd service
# Generated by Ansible - DO NOT EDIT MANUALLY
# Template: node_exporter.service.j2

[Unit]
Description=Node Exporter - PAI System Metrics Collection
Documentation=https://prometheus.io/docs/guides/node-exporter/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=node_exporter
Group=node_exporter
ExecStart=/usr/local/bin/node_exporter \
  --web.listen-address=0.0.0.0:{{ node_exporter_port }} \
  --path.procfs=/proc \
  --path.rootfs=/ \
  --path.sysfs=/sys \
  --collector.filesystem.mount-points-exclude='^/(sys|proc|dev|host|etc)($|/)' \
  --collector.systemd \
  --collector.processes \
  --no-collector.infiniband

ExecReload=/bin/kill -HUP $MAINPID
SyslogIdentifier=node_exporter
Restart=always
RestartSec=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict

# Resource limits
MemoryLimit=128M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
