# ansible/templates/pai-portal.service.j2
# PAI Portal systemd service
# Generated by Ansible - DO NOT EDIT MANUALLY
# Template: pai-portal.service.j2

[Unit]
Description=PAI Infrastructure Home Portal
Documentation=https://github.com/thebyrdman-git/personal-pai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jbyrd
Group=jbyrd
WorkingDirectory=/home/jbyrd/hatter-pai
ExecStart=/usr/bin/python3 /home/jbyrd/hatter-pai/pai-portal.py
ExecReload=/bin/kill -HUP $MAINPID
Environment=FLASK_ENV=production
Environment=FLASK_APP=pai-portal.py

SyslogIdentifier=pai-portal
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
MemoryLimit=256M
CPUQuota=50%

# Health check
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
