# ansible/templates/prometheus.service.j2
# Prometheus systemd service
# Generated by Ansible - DO NOT EDIT MANUALLY
# Template: prometheus.service.j2

[Unit]
Description=Prometheus - PAI Infrastructure Monitoring
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_user }}
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/bin/prometheus \
  --config.file={{ prometheus_config_dir }}/prometheus.yml \
  --storage.tsdb.path={{ prometheus_data_dir }} \
  --web.console.templates=/etc/prometheus/consoles \
  --web.console.libraries=/etc/prometheus/console_libraries \
  --web.listen-address=0.0.0.0:{{ prometheus_port }} \
  --web.external-url=http://{{ ansible_host }}:{{ prometheus_port }} \
  --storage.tsdb.retention.time=30d \
  --web.enable-lifecycle

SyslogIdentifier=prometheus
Restart=always
RestartSec=5
LimitNOFILE=65536

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths={{ prometheus_data_dir }}
ReadOnlyPaths={{ prometheus_config_dir }}

# Resource limits
MemoryLimit=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
